# services_core/core-fs/configure.service
[Unit]
Description=Professional Service for core-fs
After=core-fs.service
Documentation=https://github.com/professional/init-not-systemd

[Service]
Type=simple
ExecStart=./services/core-fs/run
Restart=always
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
