# services_core/swap-manager/configure.service
[Unit]
Description=Professional Service for swap-manager
After=core-fs.service
Documentation=https://github.com/professional/init-not-systemd

[Service]
Type=simple
ExecStart=./services/swap-manager/run
Restart=always
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
