# services_core/logger/configure.service
[Unit]
Description=Professional Service for logger
After=core-fs.service
Documentation=https://github.com/professional/init-not-systemd

[Service]
Type=simple
ExecStart=./services/logger/run
Restart=always
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
